# site-cookbooks/database/templates/default/mysqlbackup.cron.erb
26 13 * * * root /usr/local/bin/mysqlbackup

